# backend/launchd/com.czstore.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent para iniciar o backend CZ Store automaticamente
  sempre que você fizer login no macOS.

  Instale com: backend/launchd/install.sh
  Remova com:  backend/launchd/uninstall.sh

  Paths são substituídos pelo script de instalação.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.czstore.backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__PROJECT_PATH__/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_PATH__/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH__/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
